# boot/bunkerbox_watchdog.conf
#!upstart

description "BunkerBox Watchdog Service"
author      "Jim Heising <jheising@gmail.com>"

expect fork

start on (startup and virtual-filesystems)
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

env HOME=/bunker/control

script
    # Not sure why $HOME is needed, but we found that it is:
    export HOME=$HOME
    chdir $HOME
    exec node bunkerbox_watchdog.js >> logs/watchdog.log &
end script

#post-start script
#   # Optionally put a script here that will notifiy you node has (re)started
#   # /root/bin/hoptoad.sh "node.js has started!"
#end script